# deploy/supervisor.ini
[program:finepy]
command = /www/venv/bin/gunicorn -c /www/finepy/deploy/gunicorn.conf fine:fine
directory = /www/finepy
user = root
autostart = true
autorestart = true
