# deploy/supervisor.ini
[program:finepy]
command = /www/opt/venv/bin/gunicorn -c /www/finepy/deploy/gunicorn.conf fine:app
directory = /www/finepy
user = root 
autostart = true
autorestart = true
