# deploy/supervisor_local.ini
[program:finepy]
command = /opt/www/venv/bin/gunicorn -c /opt/www/finepy/deploy/gunicorn.conf fine:app
directory = /opt/www/finepy
user = root 
autostart = true
autorestart = true
